# .config/systemd/user/battery-notify.service
[Unit]
Description=Battery low notification

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/battery-notify.sh
