# .config/systemd/user/mpris-proxy.service
[Unit]
Description=MPRIS Proxy

[Service]
ExecStart=/usr/bin/mpris-proxy
Restart=on-failure

[Install]
WantedBy=default.target
